# dockerfiles/cron/semawi.cron
# -*- mode: crontab;-*-
0 3 * * * curl http://semawi-gc2sync:5000 >& /dev/null
