# dockerfiles/cron/semawi.cron
# -*- mode: crontab;-*-
0 */1 * * * curl -A 'semawi-cron-runjobs' http://semawi-mediawiki:5000/runjobs >& /dev/null
