# dockerfiles/cron/semawi.cron
# -*- mode: crontab;-*-
0 3 * * * curl http://semawi-gc2sync:5000 >& /dev/null
0 */1 * * * curl -A 'semawi-cron-runjobs' http://semawi-mediawiki:5000/runjobs >& /dev/null
